# docker/面向桌面应用程序的Docker容器/其他/tigervnc/arch.xfce4.tigervnc/files/etc/supervisor.d/vncserver/vncserver.ini
[program:vncserver]
priority=200
command=/usr/bin/vncserver :1 -geometry 1024x768 -fg
user=root
autostart=true
autorestart=true
stopwaitsecs=30
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_logfile_backups=0

