# deepin/ubuntu-xrdp3389设置/ubuntu_novnc/supervisor_novnc.conf
[program:xvfb]
priority=10
directory=/
command=/usr/bin/Xvfb :1 -screen 0 %(ENV_SCREEN_DIMENSIONS)s
user=root
autostart=true
autorestart=true

[program:x11vnc]
priority=20
directory=/
command=x11vnc -rfbauth /home/%(ENV_DESKTOP_USERNAME)s/.vnc/passwd -display :1 -xkb -ncache 0 -repeat -forever
user=root
autostart=true
autorestart=true

[program:novnc]
priority=25
directory=/noVNC
command=/noVNC/utils/launch.sh --vnc localhost:5900
user=root
autostart=true
autorestart=true
